# 99-midi.rules
# this file belongs into /etc/udev/rules.d/
ACTION=="add", SUBSYSTEM=="sound", RUN+="/usr/local/bin/aconnect.sh"