# timidity.service
# this file belongs into /etc/systemd/system/
# after install, run `systemctl enable timidity.service`
[Unit]
Description=TiMidity++ Daemon
After=sound.target

[Service]
ExecStart=/usr/bin/timidity -iA -Os

[Install]
WantedBy=default.target